# miAplicacion.conf
start on runlevel [5000]
stop on runlevel [!5000]

respawn
respawn limit 10 5

setuid ubuntu
chdir /var/www/sockets

env PORT=5000

exec npm start
